# etc/rc.local
#!/bin/bash
#
#
#


sudo -u eosstracker /eosstracker/bin/start_session.bash

exit 0
